# artemis/upstart/artemis-boot.conf
description "start artemis on boot"

task
start on runlevel [2345]

exec initctl emit artemis-start
