# artemis/upstart/artemis-consumers-restart.conf
description "restart queue consumers"

task

start on artemis-restart or artemis-kill

script
    . /etc/default/artemis
    manage-consumers
end script
